# toolchain/services/ibmtss-service/ibmtss.service
[Unit]
Description=IBM Virtual TPM2.0

[Service]
User=teri
Group=teri
ExecStart=/usr/bin/tpm_server

[Install]
WantedBy=multi-user.target
