# toolchain/services/ibmtss-service/ibmtss.service
[Unit]
Description=IBM Virtual TPM2.0
After=diasfw.service

[Service]
ExecStart=/usr/bin/tpm_server

[Install]
WantedBy=multi-user.target
